# systemd/flopper.service
[Unit]
After=network.target

[Service]
ExecStart=/home/pi/start.sh

[Install]
WantedBy=default.target